# 12-systemd/1/logchecker.service
[Unit]
Description=Log Checker

[Service]
Type=oneshot
EnvironmentFile=/etc/sysconfig/logchecker
ExecStart=/opt/logchecker.sh $WORD $LOG